# cron.d/sun
MAILTO=''

*/15	*	*	*	*	web22	/usr/bin/php -f /var/www/ines.cz/web/pavel/sun/sun.php sun0304 http://sdo.gsfc.nasa.gov/assets/img/latest/latest_1024_0304.jpg
*/15	*	*	*	*	web22	/usr/bin/php -f /var/www/ines.cz/web/pavel/sun/sun.php sunHMII http://sdo.gsfc.nasa.gov/assets/img/latest/latest_1024_HMII.jpg
